# debian/bolt.prerm
#!/bin/sh
set -e

case "$1" in
    remove|upgrade|deconfigure)
        # Stop the service if it's running
        if systemctl is-active --quiet bolt; then
            systemctl stop bolt || true
        fi

        # Disable the service
        if systemctl is-enabled --quiet bolt; then
            systemctl disable bolt || true
        fi
        ;;
esac

#DEBHELPER#

exit 0